# NodeA/crontab
*/3 * * * * /app/backup.sh >> /app/backup.log 2>&1
